# scripts/systemd/crm-send-prod.timer
[Unit]
Description=Run CRM email queue processor every 10 minutes (production)
Documentation=https://github.com/kb24quentin/kfzblitz24

[Timer]
OnBootSec=60s
OnUnitActiveSec=10min
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
